# systemd-user/medialist.service
[Unit]
Description=medialist update

[Service]
Type=oneshot
ExecStart=/home/jonah/doc/medialist/update

[Install]
WantedBy=timer-daily.target
