# systemd-user/timer-daily.target
[Unit]
Description=daily timer target
StopWhenUnneeded=true
